# src/karmalb-pkg/filelist
d etc
d etc/cron.d
f etc/cron.d/zenloadbalancer
f etc/firstboot
d etc/init.d
f etc/init.d/minihttpd
f etc/init.d/zenloadbalancer
d etc/logrotate.d
f etc/logrotate.d/zenloadbalancer
d etc/modprobe.d
f etc/modprobe.d/karmalb.conf
d etc/modules-load.d
f etc/modules-load.d/karmalb.conf
f etc/motd.tail
f etc/profile
d etc/rc2.d
l etc/rc2.d/S01zenloadbalancer ../init.d/zenloadbalancer
d etc/rc3.d
l etc/rc3.d/S01zenloadbalancer ../init.d/zenloadbalancer
d etc/rc4.d
l etc/rc4.d/S01zenloadbalancer ../init.d/zenloadbalancer
d etc/rc5.d
l etc/rc5.d/S01zenloadbalancer ../init.d/zenloadbalancer
d etc/sysctl.d
f etc/sysctl.d/70-karmalb.conf
l license.txt opt/klb/license.txt
d opt
d opt/klb
d opt/klb/app
d opt/klb/app/checkglobalconf
f opt/klb/app/checkglobalconf/checkglobalconf.cgi
F opt/klb/app/checkglobalconf/global.conf.tpl
f opt/klb/app/checkglobalconf/snmpd.conf.tpl
d opt/klb/app/farmguardian
d opt/klb/app/farmguardian/bin
f opt/klb/app/farmguardian/bin/farmguardian
d opt/klb/app/gdnsd
#f opt/klb/app/gdnsd/libgdnsd.so
#f opt/klb/app/gdnsd/plugin_http_status.so
#f opt/klb/app/gdnsd/plugin_multifo.so
#f opt/klb/app/gdnsd/plugin_null.so
#f opt/klb/app/gdnsd/plugin_reflect.so
#f opt/klb/app/gdnsd/plugin_simplefo.so
#f opt/klb/app/gdnsd/plugin_static.so
#f opt/klb/app/gdnsd/plugin_tcp_connect.so
d opt/klb/app/gdnsd/sbin
#f opt/klb/app/gdnsd/sbin/gdnsd
d opt/klb/app/libexec
l opt/klb/app/libexec/check_http /usr/lib/nagios/plugins/check_http
l opt/klb/app/libexec/check_ldap /usr/lib/nagios/plugins/check_ldap
l opt/klb/app/libexec/check_ldaps /usr/lib/nagios/plugins/check_ldap
l opt/klb/app/libexec/check_smtp /usr/lib/nagios/plugins/check_smtp
l opt/klb/app/libexec/check_tcp /usr/lib/nagios/plugins/check_tcp
l opt/klb/app/libexec/check_udp /usr/lib/nagios/plugins/check_udp
f opt/klb/app/libexec/check_uplink
d opt/klb/app/mini_httpd
#f opt/klb/app/mini_httpd/htpasswd
f opt/klb/app/mini_httpd/.htpasswd
l opt/klb/app/mini_httpd/mini_httpd /usr/sbin/mini_httpd
f opt/klb/app/mini_httpd/mini_httpd.conf
f opt/klb/app/mini_httpd/lbguicert.pem
d opt/klb/app/openssl
d opt/klb/app/openssl/openssl
d opt/klb/app/openssl/openssl/misc
f opt/klb/app/openssl/openssl/misc/CA.pl
f opt/klb/app/openssl/openssl/misc/CA.sh
f opt/klb/app/openssl/openssl/misc/c_hash
f opt/klb/app/openssl/openssl/misc/c_info
f opt/klb/app/openssl/openssl/misc/c_issuer
f opt/klb/app/openssl/openssl/misc/c_name
f opt/klb/app/openssl/openssl/misc/tsget
f opt/klb/app/openssl/openssl/openssl.cnf
d opt/klb/app/openssl/usr
d opt/klb/app/openssl/usr/local
d opt/klb/app/openssl/usr/local/bin
f opt/klb/app/openssl/usr/local/bin/c_rehash
#f opt/klb/app/openssl/usr/local/bin/openssl
d opt/klb/app/openssl/usr/local/lib
f opt/klb/app/openssl/usr/local/lib/libcrypto.a
f opt/klb/app/openssl/usr/local/lib/libssl.a
d opt/klb/app/openssl/usr/local/lib/pkgconfig
f opt/klb/app/openssl/usr/local/lib/pkgconfig/libcrypto.pc
f opt/klb/app/openssl/usr/local/lib/pkgconfig/libssl.pc
f opt/klb/app/openssl/usr/local/lib/pkgconfig/openssl.pc
d opt/klb/app/pound
d opt/klb/app/pound/etc
f opt/klb/app/pound/etc/poundtpl.cfg
d opt/klb/app/pound/sbin
#f opt/klb/app/pound/sbin/pound
#f opt/klb/app/pound/sbin/poundctl
d opt/klb/app/pound/share
d opt/klb/app/pound/share/man
d opt/klb/app/pound/share/man/man8
f opt/klb/app/pound/share/man/man8/pound.8
f opt/klb/app/pound/share/man/man8/poundctl.8
d opt/klb/app/ucarp
d opt/klb/app/ucarp/sbin
#f opt/klb/app/ucarp/sbin/ucarp
d opt/klb/app/zbin
f opt/klb/app/zbin/replace-nic.sh
f opt/klb/app/zbin/update-chrony.sh
f opt/klb/app/zbin/update-issue.sh
f opt/klb/app/zbin/zlb-firstboot.sh
d opt/klb/app/zenbackup
f opt/klb/app/zenbackup/diff-backup.sh
f opt/klb/app/zenbackup/extract-backup.sh
f opt/klb/app/zenbackup/zenbackup.pl
f opt/klb/app/zenbackup/zenbackup.pl
d opt/klb/app/zeninotify
f opt/klb/app/zeninotify/zeninotify.pl
d opt/klb/app/zenlatency
f opt/klb/app/zenlatency/zenlatency-start.pl
f opt/klb/app/zenlatency/zenlatency-stop.pl
d opt/klb/app/zenrrd
f opt/klb/app/zenrrd/cpu-rrd.pl
f opt/klb/app/zenrrd/disk-rrd.pl
f opt/klb/app/zenrrd/farm-rrd.pl
f opt/klb/app/zenrrd/iface-rrd.pl
f opt/klb/app/zenrrd/load-rrd.pl
f opt/klb/app/zenrrd/memory-rrd.pl
f opt/klb/app/zenrrd/memoryswap-rrd.pl
d opt/klb/app/zenrrd/rrd
f opt/klb/app/zenrrd/temperature-rrd.pl
f opt/klb/app/zenrrd/zenrrd.pl
d opt/klb/backups
l opt/klb/backups/.htpasswd ../www/.htpasswd
d opt/klb/config
F opt/klb/config/global.conf
f opt/klb/config/timeservers.conf
f opt/klb/config/lbcert.pem
f opt/klb/config/zlb-start
f opt/klb/config/zlb-stop
f opt/klb/license.txt
d opt/klb/logs
d opt/klb/www
f opt/klb/www/cert_functions.cgi
f opt/klb/www/cluster_functions.cgi
f opt/klb/www/content1-1.cgi
f opt/klb/www/content1-202.cgi
f opt/klb/www/content1-21.cgi
f opt/klb/www/content1-24.cgi
f opt/klb/www/content1-25.cgi
f opt/klb/www/content1-26.cgi
f opt/klb/www/content1-27.cgi
f opt/klb/www/content1-28.cgi
f opt/klb/www/content1-29.cgi
f opt/klb/www/content1-2.cgi
f opt/klb/www/content1-3.cgi
f opt/klb/www/content2-1.cgi
f opt/klb/www/content2-2.cgi
f opt/klb/www/content2-3.cgi
f opt/klb/www/content3-1.cgi
f opt/klb/www/content3-21.cgi
f opt/klb/www/content3-2.cgi
f opt/klb/www/content3-3.cgi
f opt/klb/www/content3-4.cgi
f opt/klb/www/content3-5.cgi
f opt/klb/www/content4-1.cgi
f opt/klb/www/content-tpl.cgi
d opt/klb/www/css
f opt/klb/www/css/base.css
f opt/klb/www/css/grid.css
f opt/klb/www/css/menu_style.css
f opt/klb/www/datalink_functions.cgi
f opt/klb/www/downloadcerts.cgi
f opt/klb/www/downloads.cgi
f opt/klb/www/farmguardian_functions.cgi
f opt/klb/www/farms_functions.cgi
d opt/klb/www/font
l opt/klb/www/font/font-aw font-awesome-4.3.0
d opt/klb/www/font/font-awesome-4.3.0
d opt/klb/www/font/font-awesome-4.3.0/css
f opt/klb/www/font/font-awesome-4.3.0/css/font-awesome.css
f opt/klb/www/font/font-awesome-4.3.0/css/font-awesome.min.css
d opt/klb/www/font/font-awesome-4.3.0/fonts
f opt/klb/www/font/font-awesome-4.3.0/fonts/FontAwesome.otf
f opt/klb/www/font/font-awesome-4.3.0/fonts/fontawesome-webfont.eot
f opt/klb/www/font/font-awesome-4.3.0/fonts/fontawesome-webfont.svg
f opt/klb/www/font/font-awesome-4.3.0/fonts/fontawesome-webfont.ttf
f opt/klb/www/font/font-awesome-4.3.0/fonts/fontawesome-webfont.woff
f opt/klb/www/font/font-awesome-4.3.0/fonts/fontawesome-webfont.woff2
d opt/klb/www/font/font-awesome-4.3.0/less
f opt/klb/www/font/font-awesome-4.3.0/less/animated.less
f opt/klb/www/font/font-awesome-4.3.0/less/bordered-pulled.less
f opt/klb/www/font/font-awesome-4.3.0/less/core.less
f opt/klb/www/font/font-awesome-4.3.0/less/fixed-width.less
f opt/klb/www/font/font-awesome-4.3.0/less/font-awesome.less
f opt/klb/www/font/font-awesome-4.3.0/less/icons.less
f opt/klb/www/font/font-awesome-4.3.0/less/larger.less
f opt/klb/www/font/font-awesome-4.3.0/less/list.less
f opt/klb/www/font/font-awesome-4.3.0/less/mixins.less
f opt/klb/www/font/font-awesome-4.3.0/less/path.less
f opt/klb/www/font/font-awesome-4.3.0/less/rotated-flipped.less
f opt/klb/www/font/font-awesome-4.3.0/less/stacked.less
f opt/klb/www/font/font-awesome-4.3.0/less/variables.less
d opt/klb/www/font/font-awesome-4.3.0/scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_animated.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_bordered-pulled.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_core.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_fixed-width.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/font-awesome.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_icons.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_larger.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_list.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_mixins.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_path.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_rotated-flipped.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_stacked.scss
f opt/klb/www/font/font-awesome-4.3.0/scss/_variables.scss
f opt/klb/www/footer.cgi
f opt/klb/www/functions.cgi
f opt/klb/www/gslb_functions.cgi
f opt/klb/www/gui_functions.cgi
f opt/klb/www/header.cgi
f opt/klb/www/help.cgi
d opt/klb/www/help
f opt/klb/www/help/clustering.html
f opt/klb/www/help-content.cgi
f opt/klb/www/.htpasswd
f opt/klb/www/http_functions.cgi
d opt/klb/www/img
f opt/klb/www/img/back.png
f opt/klb/www/img/boxheader2.gif
f opt/klb/www/img/boxheader.gif
f opt/klb/www/img/boxheader.png
f opt/klb/www/img/btn-bg-dark-active.gif
f opt/klb/www/img/btn-bg-dark-active.png
f opt/klb/www/img/btn-bg-dark-active-tp.png
f opt/klb/www/img/btn-bg-dark.gif
f opt/klb/www/img/btn-bg-dark-hover.gif
f opt/klb/www/img/favicon2.ico
f opt/klb/www/img/favicon.ico
f opt/klb/www/img/go-bg.gif
d opt/klb/www/img/graphs
f opt/klb/www/img/header-bottom.gif
f opt/klb/www/img/header-middle.gif
d opt/klb/www/img/icons
d opt/klb/www/img/icons/nav
f opt/klb/www/img/icons/nav/buttonmodify.png
f opt/klb/www/img/icons/nav/content.png
f opt/klb/www/img/icons/nav/dashboard.png
f opt/klb/www/img/icons/nav/help.png
f opt/klb/www/img/icons/nav/settings.png
f opt/klb/www/img/icons/nav/upload.png
f opt/klb/www/img/icons/nav/users.png
d opt/klb/www/img/icons/notifications
f opt/klb/www/img/icons/notifications/error.png
f opt/klb/www/img/icons/notifications/info.png
f opt/klb/www/img/icons/notifications/success.png
f opt/klb/www/img/icons/notifications/tip.png
f opt/klb/www/img/icons/notifications/warning.png
d opt/klb/www/img/icons/small
f opt/klb/www/img/icons/small/accept2.png
f opt/klb/www/img/icons/small/accept.png
f opt/klb/www/img/icons/small/add.png
f opt/klb/www/img/icons/small/application_edit.png
f opt/klb/www/img/icons/small/application_home.png
f opt/klb/www/img/icons/small/arrow_down.png
f opt/klb/www/img/icons/small/arrow_left.png
f opt/klb/www/img/icons/small/arrow_refresh.png
f opt/klb/www/img/icons/small/arrow_up.png
f opt/klb/www/img/icons/small/bullet_toggle_minus.png
f opt/klb/www/img/icons/small/bullet_toggle_plus.png
f opt/klb/www/img/icons/small/cancel.png
f opt/klb/www/img/icons/small/cart_put.png
f opt/klb/www/img/icons/small/chart_bar.png
f opt/klb/www/img/icons/small/close.png
f opt/klb/www/img/icons/small/cog_add.png
f opt/klb/www/img/icons/small/cog_delete.png
f opt/klb/www/img/icons/small/cog_edit.png
f opt/klb/www/img/icons/small/cog_error.png
f opt/klb/www/img/icons/small/cog_go.png
f opt/klb/www/img/icons/small/cog_return.png
f opt/klb/www/img/icons/small/cog_save.png
f opt/klb/www/img/icons/small/comment_add.png
f opt/klb/www/img/icons/small/comment_delete.png
f opt/klb/www/img/icons/small/comment_edit.png
f opt/klb/www/img/icons/small/comment.png
f opt/klb/www/img/icons/small/connect.png
f opt/klb/www/img/icons/small/cross_octagon.png
f opt/klb/www/img/icons/small/cross.png
f opt/klb/www/img/icons/small/database.png
f opt/klb/www/img/icons/small/databases.png
f opt/klb/www/img/icons/small/delete.png
f opt/klb/www/img/icons/small/disconnect.png
f opt/klb/www/img/icons/small/disk.png
f opt/klb/www/img/icons/small/drive_disk.png
f opt/klb/www/img/icons/small/exclamation_octagon_fram.png
f opt/klb/www/img/icons/small/exclamation.png
f opt/klb/www/img/icons/small/farm_add.png
f opt/klb/www/img/icons/small/farm_cancel.png
f opt/klb/www/img/icons/small/farm_delete.png
f opt/klb/www/img/icons/small/farm_edit.png
f opt/klb/www/img/icons/small/farm_manage.png
f opt/klb/www/img/icons/small/farm.png
f opt/klb/www/img/icons/small/farms.png
f opt/klb/www/img/icons/small/farm_up.png
f opt/klb/www/img/icons/small/find.png
f opt/klb/www/img/icons/small/help.png
f opt/klb/www/img/icons/small/image_add.png
f opt/klb/www/img/icons/small/image_delete.png
f opt/klb/www/img/icons/small/image_edit.png
f opt/klb/www/img/icons/small/image_manage.png
f opt/klb/www/img/icons/small/images.png
f opt/klb/www/img/icons/small/information.png
f opt/klb/www/img/icons/small/key.png
f opt/klb/www/img/icons/small/lock.png
f opt/klb/www/img/icons/small/page_2.png
f opt/klb/www/img/icons/small/page_add.png
f opt/klb/www/img/icons/small/page_delete.png
f opt/klb/www/img/icons/small/page_edit.png
f opt/klb/www/img/icons/small/page_lightning.png
f opt/klb/www/img/icons/small/page.png
f opt/klb/www/img/icons/small/page_white_add.png
f opt/klb/www/img/icons/small/page_white_delete.png
f opt/klb/www/img/icons/small/page_white_find.png
f opt/klb/www/img/icons/small/page_white_get.png
f opt/klb/www/img/icons/small/page_white_put.png
f opt/klb/www/img/icons/small/pencil.png
f opt/klb/www/img/icons/small/plugin_add.png
f opt/klb/www/img/icons/small/plugin_back.png
f opt/klb/www/img/icons/small/plugin_delete2.png
f opt/klb/www/img/icons/small/plugin_delete.png
f opt/klb/www/img/icons/small/plugin_disabled.png
f opt/klb/www/img/icons/small/plugin_edit.png
f opt/klb/www/img/icons/small/plugin_error.png
f opt/klb/www/img/icons/small/plugin_go.png
f opt/klb/www/img/icons/small/plugin_link.png
f opt/klb/www/img/icons/small/plugin.png
f opt/klb/www/img/icons/small/plugin_save.png
f opt/klb/www/img/icons/small/plugin_stop.png
f opt/klb/www/img/icons/small/plugin_upn.png
f opt/klb/www/img/icons/small/pluginv_add.png
f opt/klb/www/img/icons/small/rosette.png
f opt/klb/www/img/icons/small/script.png
f opt/klb/www/img/icons/small/server_add.png
f opt/klb/www/img/icons/small/server_database.png
f opt/klb/www/img/icons/small/server_delete.png
f opt/klb/www/img/icons/small/server_edit.png
f opt/klb/www/img/icons/small/server_maintenance.png
f opt/klb/www/img/icons/small/server_ok.png
f opt/klb/www/img/icons/small/server_out.png
f opt/klb/www/img/icons/small/server_save.png
f opt/klb/www/img/icons/small/start.png
f opt/klb/www/img/icons/small/stop.png
f opt/klb/www/img/icons/small/table_multiple.png
f opt/klb/www/img/icons/small/toggle_close.png
f opt/klb/www/img/icons/small/toggle_open.png
f opt/klb/www/img/icons/small/user_add.png
f opt/klb/www/img/icons/small/user_delete.png
f opt/klb/www/img/icons/small/user_edit.png
f opt/klb/www/img/icons/small/user.png
f opt/klb/www/img/icons/small/warning.png
f opt/klb/www/img/icons/small/zoom_in.png
f opt/klb/www/img/img_verified_logo.gif
f opt/klb/www/img/input-bg.gif
f opt/klb/www/img/logo.png
d opt/klb/www/img/menu_style
f opt/klb/www/img/menu_style/hover.gif
f opt/klb/www/img/menu_style/hover_sub.gif
f opt/klb/www/img/menu_style/seperator.gif
f opt/klb/www/img/menu_style/sub_sep.gif
f opt/klb/www/img/page-bg.gif
f opt/klb/www/img/search-bg.gif
f opt/klb/www/img/sidebar-bg.png
f opt/klb/www/index.cgi
f opt/klb/www/l4_functions.cgi
f opt/klb/www/menu.cgi
f opt/klb/www/networking_functions.cgi
f opt/klb/www/nf_functions.cgi
d opt/klb/www/out
d opt/klb/www/out/css
f opt/klb/www/out/css/base.css
d opt/klb/www/out/img
f opt/klb/www/out/img/back.png
f opt/klb/www/out/img/boxheader2.gif
f opt/klb/www/out/img/boxheader.gif
f opt/klb/www/out/img/boxheader.png
f opt/klb/www/out/img/btn-bg-dark-active.gif
f opt/klb/www/out/img/btn-bg-dark-active.png
f opt/klb/www/out/img/btn-bg-dark-active-tp.png
f opt/klb/www/out/img/btn-bg-dark.gif
f opt/klb/www/out/img/btn-bg-dark-hover.gif
f opt/klb/www/out/img/favicon2.ico
f opt/klb/www/out/img/favicon.ico
f opt/klb/www/out/img/go-bg.gif
f opt/klb/www/out/img/header-bottom.gif
f opt/klb/www/out/img/header-middle.gif
f opt/klb/www/out/img/input-bg.gif
f opt/klb/www/out/img/logo.png
f opt/klb/www/out/img/page-bg.gif
f opt/klb/www/out/img/search-bg.gif
f opt/klb/www/out/img/sidebar-bg.png
f opt/klb/www/out/out.cgi
f opt/klb/www/out/out.html
f opt/klb/www/rrd_functions.cgi
f opt/klb/www/snmp_functions.cgi
f opt/klb/www/system_functions.cgi
f opt/klb/www/uploadcerts.cgi
f opt/klb/www/uploadcertsfromcsr.cgi
f opt/klb/www/upload.cgi
