# src/karmalb-pkg/files/etc/cron.d/zenloadbalancer
*/5 * * * *	root	/opt/klb/app/zenrrd/zenrrd.pl &>/dev/null
