# src/karmalb-pkg/files/etc/cron.d/zenloadbalancer
*/5 * * * *	root	/opt/klb/app/zenrrd/zenrrd.pl &>/dev/null
*/10  * * * * 	root	/opt/klb/app/zenntp/zenntp.pl &>/dev/null
