# src/karmalb-pkg/control/postinst
#!/bin/bash


[ "$1" == "configure" ] && echo "Concluding ZLB installation..." 

update-rc.d zenloadbalancer defaults

exit 0
