# src/karmalb-pkg/control/postinst
#!/bin/bash


[ "$1" == "configure" ] && echo "Concluding @SHORTNAME@ installation..." 

update-rc.d zenloadbalancer defaults

sed -i -e '1i\
@LONGNAME@ @VERSION@' /etc/issue

exit 0
